# systemd/cloudflare-sync.service
[Unit]
Description=Cloudflare Registry Sync Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/cloudflare-sync/tools
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/opt/cloudflare-sync/tools"
ExecStart=/usr/bin/python3 /opt/cloudflare-sync/tools/webhook_listener.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/cloudflare-sync/service.log
StandardError=append:/var/log/cloudflare-sync/error.log

[Install]
WantedBy=multi-user.target